# bonk.service
[Unit]
Description=Run Bonk IP ban hammer

[Service]
ExecStart=/usr/bin/bonk --is-agent
User=root

[Install]
WantedBy=multi-user.target